# systemd/user/kiosk-agent.service
# kiosk-client agent user service
#
# This service starts the local kiosk-agent inside the graphical user session.
# The agent exposes local HTTP status endpoints and monitors the browser
# runtime without controlling it.

[Unit]
Description=kiosk-client agent
After=graphical-session.target
BindsTo=graphical-session.target
PartOf=graphical-session.target

[Service]
WorkingDirectory=%h/kiosk-client
ExecStart=%h/kiosk-client/kiosk-agent
Restart=always
RestartSec=3

[Install]
WantedBy=graphical-session.target
